# dev-with-docker/code-server/etc/services.d/code-server/run
#!/usr/bin/with-contenv sh

exec code-server --config /app/code-server/config/config.yaml